# backend/deploy/systemd/resilience-backend.service.example
[Unit]
Description=Resilience Training Programme Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/resilience-training-programme/backend
EnvironmentFile=/opt/resilience-training-programme/backend/.env
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
